# build/boxes/ElasticKube/files/kubegrunt.conf
#!upstart
description "ElasticKube Grunt"

setuid vagrant

start on runlevel [2345]
stop  on shutdown

respawn
respawn limit unlimited

script
     # Run npm install
    cd /opt/elastickube/src/ui
    npm install

    cd /opt/elastickube/docs
    npm install

    cd /opt/elastickube/build/kubegrunt
    npm install

    # Start services
    cd /opt/elastickube/src/ui
    npm start > /var/log/elastickube/ui.log &

    # Build documentation
    export PYTHONPATH=/opt/elastickube/src
    cd /opt/elastickube/docs
    npm start > /var/log/elastickube/documentation.log &

    # Launch Grunt
    cd /opt/elastickube/build/kubegrunt
    grunt > /var/log/elastickube/build.log &

    wait
end script

