# build/boxes/ElasticKube Dev/files/kubegrunt.conf
#!upstart
description "ElasticKube Grunt"

setuid elasticbox

start on runlevel [2345]
stop  on shutdown

respawn
respawn limit unlimited

script
    sudo mkdir -p /var/log/elastickube
    sudo chown elasticbox:elasticbox /var/log/elastickube

    cd /opt/elastickube/src/ui
    npm install
    npm start > /var/log/elastickube/ui.log &

    # Launch Grunt
    cd /opt/elastickube/build/kubegrunt
    npm install
    grunt > /var/log/elastickube/build.log &

    wait
end script

